# docker/api/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

#[program:nodecollector]
#command=nodemon /var/user1/node_collector &